# crontabs/accDailySummary.d0en
5 0 * * * enstore /home/enstore/enstore/sbin/ecron -p accDailySummary -l /home/enstore/enstore/src/acc_daily_summary.py > /dev/null 2> /dev/null
