# crontabs/accounting_db
#30 4 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p vacuumACC -l db_vacuum.py accounting > /dev/null 2> /dev/null
#30 4 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p vacuumACC -l vacuumdb -p 8800 -d accounting -v -z > /dev/null 2> /dev/null
30 2 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p pg_base_backup_accounting /opt/enstore/tools/pg_base_backup.sh accounting >/dev/null 2>&1
