# crontabs/aml2mirror.d0en
18 9 * * * enstore /home/enstore/enstore/sbin/ecron -p aml2mirror mirror /home/enstore/CRON/aml2.mirror > /dev/null 2> /dev/null
