# crontabs/aml2logs.d0ensrv2
#The following make the AML2 log web pages.
10,25,40,55 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p aml2logs aml2logs > /dev/null 2> /dev/null
0 4 15 * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p aml2log-stash aml2log-stash > /dev/null 2> /dev/null
