# crontabs/checkPNFS
3,13,23,33,43,53 * * * * root /home/enstore/enstore/sbin/ecron -p checkPNFS checkPNFS >/dev/null 2>/dev/null
