# crontabs/checkPNFS
3,13,23,33,43,53 * * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p checkPNFS checkPNFS >/dev/null 2>/dev/null
