# crontabs/chkcrons
50 * * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p chkcrons nice -n 20 chkcrons  >/dev/null 2>/dev/null
