# crontabs/chkcrons
50 * * * * root /home/enstore/enstore/sbin/ecron -p chkcrons nice -n 20 chkcrons  >/dev/null 2>/dev/null
