# crontabs/cleaning_report.stken
30 0 * * * enstore /home/enstore/enstore/sbin/ecron -p cleaning_report '/usr/bin/rsh fntt -l acsss "(echo query clean all; echo logoff) | bin/cmd_proc -l -q 2>/dev/null" | tee /tmp/CLEANING | mail -s "FNTT (STKen+D0en) Drive Cleaning Summary" enstore-admin@fnal.gov; . /home/enstore/gettkt; enrcp /tmp/CLEANING stkensrv2:/fnal/ups/prd/www_pages/enstore/tape_inventory/CLEANING' >/dev/null 2>/dev/null
31 0 * * * enstore /home/enstore/enstore/sbin/ecron -p cleaning_report_GCC '/usr/bin/rsh fntt-gcc -l acsss "(echo query clean all; echo logoff) | bin/cmd_proc -l -q 2>/dev/null" | tee /tmp/CLEANING1 | mail -s "FNTT-GCC (STKen GCC) Drive Cleaning Summary" enstore-admin@fnal.gov; . /home/enstore/gettkt; enrcp /tmp/CLEANING1 stkensrv2:/fnal/ups/prd/www_pages/enstore/tape_inventory/CLEANING1' >/dev/null 2>/dev/null
