# crontabs/cleaning_report.cdfen
30 0 * * * enstore /home/enstore/enstore/sbin/ecron -p cleaning_report '/usr/bin/rsh fntt2 -l acsss "(echo query clean all; echo logoff) | bin/cmd_proc -l -q 2>/dev/null" | tee /tmp/CLEANING | mail -s "FNTT2 (CDFen) Drive Cleaning Summary" enstore-admin@fnal.gov; . /home/enstore/gettkt; enrcp /tmp/CLEANING cdfensrv2:/fnal/ups/prd/www_pages/enstore/tape_inventory/CLEANING' >/dev/null 2>/dev/null
