# crontabs/copy_ran_file
27 * * * * root /home/enstore/enstore/sbin/ecron -p copy_ran_file /home/enstore/enstore/sbin/copy_ran_file 1 alsofirstlast > /dev/null 2> /dev/null
