# crontabs/copy_ran_file.d0en
27 * * * * enstore /home/enstore/enstore/sbin/ecron -p copy_ran_file /home/enstore/enstore/sbin/choose_ran_file -C 1 \"shelf\" \"samlto\|samlto2\|mezsilo\|D0-9940B\" alsofirstlast \"full\|readonly\" > /dev/null 2> /dev/null
