# crontabs/delfile
0-55/5 * * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -q -p delfile $ENSTORE_DIR/sbin/delfile.py >/dev/null 2>/dev/null
