# crontabs/drives_info
48 06 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p drives_info drives_info.py  > /dev/null 2> /dev/null
