# crontabs/DrvBusy2.d0en
0-59/15 * * * * enstore /home/enstore/enstore/sbin/ecron -p ADICDrvBusy /home/enstore/enstore/sbin/ADICDrvBusy get >/dev/null 2>/dev/null
0-59/15 * * * * enstore /home/enstore/enstore/sbin/ecron -p STKDrvBusy  /home/enstore/enstore/sbin/STKDrvBusy  get >/dev/null 2>/dev/null
