# crontabs/DrvBusy.cdfen
5 * * * * enstore /home/enstore/enstore/sbin/ecron -p STKDrvBusy  /home/enstore/enstore/sbin/STKDrvBusy plot > /dev/null 2> /dev/null
