# crontabs/enstore_misc.stken
15 * * * * root sh -c 'chown  -R enstore.enstore /home/enstore/enstore; cd /home/enstore/enstore; chmod g+w etc etc/*.conf etc/CVS etc/CVS/Entries*'
17 00 * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p cppnfsFastBackupHISTOGRAM enrcp stkensrv1:/root/CRON/pnfsFastBackupHISTOGRAM /tmp >/dev/null 2>/dev/null

#03 00 * * 1 root /local/ups/prd/apache/v1_3_26/ups/monthlystats

#10 9,14,19 * * * root /home/enstore/enstore/sbin/enrcp d0ensrv2:/diska/tape-inventory/VOLUMES_DEFINED /tmp/d0-VOLUMES_DEFINED >/dev/null 2>/dev/null

#10 * * * * enstore /home/joes/vc/production/aaacron.ksh > /dev/null 2> /dev/null
#01 10 1-7 * 2 enstore /home/enstore/.bashrc ; /fnal/ups/prd/python/v2_2_3_E_1/bin/python /home/enstore/isa-tools/bin/RedBallAlarm.py > /tmp/redballalarm.log 2> /tmp/redballalarm.log
#0 0 * * * enstore /home/enstore/enstore/sbin/ecron -p check_multiple /home/enstore/check_multiple > /dev/null 2> /dev/null

#5 0 * * * enstore /home/enstore/enstore/sbin/ecron -p accDailySummary -l /home/enstore/enstore/src/acc_daily_summary.py > /dev/null 2> /dev/null
#0,20,40 * * * * enstore /home/enstore/enstore/sbin/ecron -p filler "sh -c 'cd /home/enstore; . /usr/local/etc/setups.sh 1>>/dev/null 2>&1; setup postgres; python /home/enstore/enstore/sbin/filler.py'"  > /dev/null 2> /dev/null
