# crontabs/enstore_html
# generate main enstore html page
55 * * * * enstore  source /usr/local/etc/setups.sh; setup enstore; python $ENSTORE_DIR/src/enstore_system_html.py > /dev/null 2> /dev/null

# gererate additional enstore web pages.
16 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p make_quota_plot_page make_quota_plot_page  > /dev/null 2> /dev/null
10 * * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p make_cron_plot_page make_cron_plot_page > /dev/null 2> /dev/null
#This last one used to only run on STKen???
35 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p make_tape_ingest_page $ENSTORE_DIR/src/make_ingest_rates_html_page.py > /dev/null 2> /dev/null

# generate SAAG
*/5 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -no-enstore -p enstoreSystem enstore system > /dev/null 2> /dev/null

# generate Active Network log page
*/5 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p enstoreNetwork enstore network > /dev/null 2> /dev/null

# generate file that contains bytes
50 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron  -p user_bytes python $ENSTORE_DIR/src/get_total_bytes_counter.py  > /dev/null 2> /dev/null

# make the web pages that show the generated plots and pages
*/5 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p make_plot_page  python $ENSTORE_DIR/src/html_main.py --plots --pages > /dev/null 2> /dev/null
