# crontabs/enstore_db
10 * * * * enstore /home/enstore/enstore/sbin/ecron -p backup enstore backup > /dev/null 2> /dev/null
30 2 * * * enstore /home/enstore/enstore/sbin/ecron -p vacuumdb -l db_vacuum.py enstoredb > /dev/null 2> /dev/null
#30 2 * * * enstore /home/enstore/enstore/sbin/ecron -p vacuumdb -l vacuumdb -p 8888 -d enstoredb -v -z > /dev/null 2> /dev/null
