# crontabs/enstore_plots
#The crontabs with two pound signs at the beginning of the comment are
# not known to be SDE ready.

# make standard plots
30 1,6,15 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p inqPlotUpdate -setup gnuplot python $ENSTORE_DIR/src/plotter.py --encp --mount > /dev/null 2> /dev/null
# storage group activity - instead of make_sg_plot
50 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; cd /tmp; $ENSTORE_DIR/sbin/ecron -p sgPlotUpdate -setup gnuplot $ENSTORE_DIR/sbin/make_sg_plot > /dev/null 2> /dev/null
#
##30 2,7,16 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p make_overall_plot python $ENSTORE_DIR/src/plotter.py --total_bytes --pts_nodes=d0ensrv2,stkensrv2,cdfensrv2 --no-plot-html  > /dev/null 2> /dev/null
30 2,7,16 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p make_overall_plot  python $ENSTORE_DIR/src/plotter_main.py -b > /dev/null 2> /dev/null


# mounts plots
30 1,6,15 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p mounts_plot -setup gnuplot python $ENSTORE_DIR/src/plotter_main.py -m > /dev/null 2> /dev/null
# ratekeeper plot
3,13,33,48 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p makeplot -setup gnuplot python $ENSTORE_DIR/src/plotter_main.py -r > /dev/null 2> /dev/null
#drive utilization
3,13,33,48 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p driveUtilization -setup gnuplot python $ENSTORE_DIR/src/plotter_main.py --utilization > /dev/null 2> /dev/null
# drive usage hours (stacked)
55 0 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p drivehours -setup gnuplot python $ENSTORE_DIR/src/plotter_main.py --drive-hours > /dev/null 2> /dev/null
# drive usage hours (separated) (same schedule as stacked)
55 0 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p drivehours_sep -setup gnuplot python $ENSTORE_DIR/src/plotter_main.py --drive-hours-sep > /dev/null 2> /dev/null
#slot utilization
10 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; cd /tmp; $ENSTORE_DIR/sbin/ecron -p slotUtilization -setup gnuplot python $ENSTORE_DIR/src/plotter_main.py --slots > /dev/null 2> /dev/null
# encp rate multi
25 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p encp_rate_plot_multi $ENSTORE_DIR/src/plotter_main.py -e >/dev/null 2> /dev/null
# file family usage
50 23 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p file_family_usage $ENSTORE_DIR/src/plotter_main.py -f > /dev/null 2> /dev/null
# files read and written per mount (stacked)
45 0 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p filesrw -setup gnuplot python $ENSTORE_DIR/src/plotter_main.py --files-rw > /dev/null 2> /dev/null
# files read and written per mount (separated) (same schedule as stacked)
45 0 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p filesrw_sep -setup gnuplot python $ENSTORE_DIR/src/plotter_main.py --files-rw-sep > /dev/null 2> /dev/null
# quotas
5 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p quotas  $ENSTORE_DIR/src/plotter_main.py -q > /dev/null 2> /dev/null
# pnfs backup time
15 1,13 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p pnfs_backup_time $ENSTORE_DIR/src/plotter_main.py -p > /dev/null 2> /dev/null
# tape burn rates
0 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p burn-rate  $ENSTORE_DIR/src/plotter_main.py -t > /dev/null 2> /dev/null
#migration summary
41 7 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p migration_summary_plots $ENSTORE_DIR/src/plotter_main.py --migration-summary > /dev/null 2> /dev/null
#mover summary
30 1 * * 1 enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p mover_summary_plots $ENSTORE_DIR/src/plotter_main.py -M > /dev/null 2> /dev/null
# fills slot usage in DB
30 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p fill_slot_usage python $ENSTORE_DIR/src/fill_slot_usage.py > /dev/null 2> /dev/null
# plots mount latency plots
30 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p mount_latencies $ENSTORE_DIR/src/plotter_main.py -l > /dev/null 2> /dev/null
# plots mount/day per tape library
30 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p library_mounts $ENSTORE_DIR/src/plotter_main.py -L > /dev/null 2> /dev/null
# Small Files Aggregation Statistics
35 1 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p sfa_stats $ENSTORE_DIR/src/plotter_main.py -S > /dev/null 2> /dev/null
