# crontabs/enstore_plots
#The crontabs with two pound signs at the beginning of the comment are
# not known to be SDE ready.

# make standard plots
30 1 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p inqPlotUpdate -setup gnuplot python $ENSTORE_DIR/src/plotter.py --encp --mount > /dev/null 2> /dev/null
# mounts plots
30 1 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p mounts_plot -setup gnuplot python $ENSTORE_DIR/src/plotter_main.py -m > /dev/null 2> /dev/null
# ratekeeper plot
3,13,33,48 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p makeplot -setup gnuplot python $ENSTORE_DIR/src/plotter_main.py -r > /dev/null 2> /dev/null
# storage grpup activity - instead of make_sg_plot
50 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p sgPlotUpdate -setup gnuplot $ENSTORE_DIR/sbin/make_sg_plot > /dev/null 2> /dev/null
25 * * * * enstore  source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p encp_rate_plot_multi $ENSTORE_DIR/sbin/encp_rate_plot_multi.py >/dev/null 2> /dev/null
50 23 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p file_family_usage $ENSTORE_DIR/sbin/file_family_analysis.py > /dev/null 2> /dev/null


##plot the quotas
5 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p quotas  $ENSTORE_DIR/sbin/quotas > /dev/null 2> /dev/null

##plot the pnfs backup time
26 00 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p pnfs_backup_time $ENSTORE_DIR/src/pnfs_backup_plot.py > /dev/null 2> /dev/null


##
25 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p encp_rate_plot_multi "sh -c 'cd /home/enstore; source /usr/local/etc/setups.sh 1>>/dev/null 2>&1; setup postgres; python $ENSTORE_DIR/sbin/encp_rate_plot_multi.py'"  > /dev/null 2> /dev/null
##
40 17 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/ecron -p file_family_usage "sh -c 'cd /home/enstore; . /usr/local/etc/setups.sh 1>>/dev/null 2>&1; setup postgres; python $ENSTORE_DIR/sbin/file_family_analysis.py'"  > /dev/null 2> /dev/null
##
# 30 2 * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p make_overall_plot python $ENSTORE_DIR/src/plotter.py --total_bytes --pts_nodes=d0ensrv2,stkensrv2,cdfensrv2 --pts_dir=/local/ups/prd/www_pages/enstore --no-plot-html  > /dev/null 2> /dev/null
30 2 * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p make_overall_plot python $ENSTORE_DIR/src/plotter.py --total_bytes --pts_nodes=d0ensrv2,stkensrv2,cdfensrv2 --no-plot-html  > /dev/null 2> /dev/null


##
20 1 * * 1 enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p weekly_summary -l $ENSTORE_DIR/src/weekly_summary_report.py /tmp/weekly_report > /dev/null 2> /dev/null

