# crontabs/enstore_html
# generate main enstore html page
55 * * * * enstore  . /usr/local/etc/setups.sh; setup enstore; python $ENSTORE_DIR/src/enstore_system_html.py > /dev/null 2> /dev/null

# generate SAAG
59/5 * * * * enstore . /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p enstoreSystem enstore system > /dev/null 2> /dev/null

# generate Active Network log page
59/5 * * * * enstore . /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p enstoreNetwork enstore network > /dev/null 2> /dev/null
