# crontabs/enstore_migration.stken
35 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p make_tape_ingest_page $ENSTORE_DIR/src/make_ingest_rates_html_page.py > /dev/null 2> /dev/null
30 10,12,16,22 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/isa-tools/bin/migration-summary > /diska/tape-inventory/MIGRATION_STATUS 2> /dev/null
