# crontabs/inventory.cms
40 * * * * enstore /home/enstore/enstore/sbin/ecron -p recent_file_listing_cms /home/enstore/enstore/src/recent_file_listing.py cms 12 /diska/tape-inventory/RECENT_FILES_ > /dev/null 2> /dev/null
