# crontabs/inventory_web
31 0 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/ecron -p cleaning_report $ENSTORE_DIR/sbin/cleaning_report > /dev/null 2> /dev/null
0-59/15 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p noaccess $ENSTORE_DIR/sbin/noaccess-tapes > /dev/null 2> /dev/null
30 1,7,13,19 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p robot_inventory $ENSTORE_DIR/sbin/Vols > /dev/null 2> /dev/null
30 8,16 * * * enstore  source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p quota_alert $ENSTORE_DIR/sbin/quota_alert > /dev/null 2> /dev/null
