# crontabs/inventory_web
#The crontabs with two pound signs at the beginning of the comment are
# not known to be SDE ready.

31 0 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p cleaning_report $ENSTORE_DIR/sbin/cleaning_report > /dev/null 2> /dev/null
0-59/15 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p noaccess $ENSTORE_DIR/sbin/noaccess-tapes > /dev/null 2> /dev/null
30 1,7,13,19 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p robot_inventory $ENSTORE_DIR/sbin/Vols > /dev/null 2> /dev/null
30 8,16 * * * enstore  source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p quota_alert $ENSTORE_DIR/sbin/quota_alert > /dev/null 2> /dev/null

1 7 * * * enstore  source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p migration_summary $ENSTORE_DIR/sbin/migration_summary > /dev/null 2> /dev/null
11 7 * * * enstore  source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p duplication_summary $ENSTORE_DIR/sbin/duplication_summary > /dev/null 2> /dev/null
45 2 * * * enstore  source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p multiple_copies_summary $ENSTORE_DIR/sbin/multiple_copy_summary > /dev/null 2> /dev/null
31 7 * * * enstore  source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p cloning_summary $ENSTORE_DIR/sbin/cloning_summary > /dev/null 2> /dev/null
41 7 * * * enstore  source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p sparse_volumes $ENSTORE_DIR/sbin/sparse_volumes > /dev/null 2> /dev/null


20 1 * * 1 enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p weekly_summary -l $ENSTORE_DIR/sbin/weekly_summary_report > /dev/null 2> /dev/null
