# crontabs/log_server
0 4 15 * * enstore /home/enstore/enstore/sbin/ecron -p log-stash log-stash > /dev/null 2> /dev/null
01 * * * * enstore /home/enstore/enstore/sbin/ecron -p traceback_check /home/enstore/enstore/sbin/check_for_traceback > /dev/null 2> /dev/null
5,20,35,50 * * * * enstore /home/enstore/enstore/sbin/ecron -p rdist-log rdist-log  > /dev/null 2> /dev/null
