# crontabs/log_server.fermi.srv2
#rdist-log is Fermilab specific.
5,20,35,50 * * * * enstore /home/enstore/enstore/sbin/ecron -p rdist-log rdist-log  > /dev/null 2> /dev/null
