# crontabs/log_html.d0en
11,41 * * * * root /home/enstore/enstore/sbin/ecron -p getnodeinfo getnodeinfo  /fnal/ups/prd/www_pages/enstore/nodeinfo.html  1>/dev/null 2>/dev/null
9,39 * * * * root /home/enstore/enstore/sbin/ecron -p udpclog udp_clog_check /fnal/ups/prd/www_pages/enstore/UDPclog.html  1>/dev/null 2>/dev/null
31 * * * * enstore /home/enstore/enstore/sbin/ecron -p PNFSRATE PNFSRATE  > /dev/null 2> /dev/null
7,37 * * * * root /home/enstore/enstore/sbin/ecron -p getcons getcons  1>/dev/null 2>/dev/null
10,25,40,55 * * * * enstore /home/enstore/enstore/sbin/ecron -p aml2logs aml2logs > /dev/null 2> /dev/null
45 * * * * enstore /home/enstore/enstore/sbin/ecron -p failedX  '/home/enstore/enstore/src/log_trans_fail.py month | tee /diska/enstore-log/transfer_failed.txt' >/dev/null 2> /dev/null
