# crontabs/log_html
#These are common.
11,41 * * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p getnodeinfo getnodeinfo  1> /dev/null 2> /dev/null
45 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p failedX $ENSTORE_DIR/src/log_trans_fail.py month > /dev/null 2> /dev/null
#45 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p failedX log_trans_fail > /dev/null 2> /dev/null
#9,39 * * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p udpclog udp_clog_check /fnal/ups/prd/www_pages/enstore/UDPclog.html  1>/dev/null 2>/dev/null
#31 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p PNFSRATE PNFSRATE  > /dev/null 2> /dev/null

#These are Fermilab specific.
#7,37 * * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p getcons getcons  1>/dev/null 2>/dev/null

#The following make the STK log web pages.
02,17,32,47 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p STKlog STKlog > /dev/null 2> /dev/null
