# crontabs/plots.stken
# added by Dmitry:
0,20,40 * * * * enstore /home/enstore/enstore/sbin/ecron -p filler "sh -c 'cd /home/enstore; . /usr/local/etc/setups.sh 1>>/dev/null 2>&1; setup postgres; python /home/enstore/enstore/sbin/filler.py'"  > /dev/null 2> /dev/null
#12 20 * * * enstore /home/enstore/enstore/sbin/ecron -p write-tabs "sh -c 'cd /home/enstore/write_tabs; . /usr/local/etc/setups.sh 1>>/dev/null 2>&1; setup postgres; /home/enstore/enstore/sbin/plot_write_tabs.py'"  > /dev/null 2> /dev/null
# by Dmitry 10/10/05:
30 1 * * * enstore /home/enstore/enstore/sbin/ecron -p write-tabs "sh -c 'cd /home/enstore/write_tabs; . /usr/local/etc/setups.sh 1>>/dev/null 2>&1; setup postgres; /home/enstore/enstore/sbin/write_tabs.py'"  > /dev/null 2> /dev/null
30 1 * * * enstore /home/enstore/enstore/sbin/ecron -p xfer-stat-daily "sh -c 'cd /home/enstore/xfer_stat; . /usr/local/etc/setups.sh 1>>/dev/null 2>&1; setup postgres; /home/enstore/enstore/sbin/xfer_stats_daily.py'"  > /dev/null 2> /dev/null
00 2 * * * enstore /home/enstore/enstore/sbin/ecron -p xfer-stat-monthly "sh -c 'cd /home/enstore/xfer_stat; . /usr/local/etc/setups.sh 1>>/dev/null 2>&1; setup postgres; /home/enstore/enstore/sbin/xfer_stats_monthly.py'"  > /dev/null 2> /dev/null
20 1 * * 1 enstore /home/enstore/enstore/sbin/ecron -p weekly_summary -l /home/enstore/enstore/src/weekly_summary_report.py /tmp/weekly_report > /dev/null 2> /dev/null
