# crontabs/pnfs_misc
10,40 * * * * root /home/enstore/enstore/sbin/ecron -p plog -setup pnfs plog >/dev/null 2>/dev/null
12 * * * * root /home/enstore/enstore/sbin/ecron -p pnfsExports -setup pnfs pnfsExportList >/dev/null 2>/dev/null
