# crontabs/pnfsFastBackup.stken
###VP 15 0-23/6 * * * root /home/enstore/enstore/sbin/ecron -p pnfsFastBackup -setup pnfs pnfsFastBackup >/dev/null 2>/dev/null
###VP 25 0-23/6 * * * root /home/enstore/enstore/sbin/ecron -p pnfsFastBackup -setup "pnfs -t" psqlBackup.1 >/dev/null 2>/dev/null
15 0,12 * * * root /home/enstore/enstore/sbin/ecron -p pnfsFastBackup -setup pnfs nice -n 20 psqlBackupV8  >/dev/null 2>/dev/null
#15 6,20 * * * root /home/enstore/enstore/sbin/ecron -p pnfsFastBackup -setup pnfs pnfsFastBackup >/dev/null 2>/dev/null
