# crontabs/pnfs_misc
#This cronjob creates and enstore log web page containing all the nodes
# that can mount pnfs.
12 * * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p pnfsExports -setup pnfs pnfsExportList >/dev/null 2>/dev/null
12 5 * * * root source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p pnfs_monitor -setup pnfs pnfs_monitor >/dev/null 2>/dev/null
