# crontabs/robot_inventory.cdfen
30 1,7,13,19 * * * enstore /home/enstore/enstore/sbin/ecron -p STKrobot_inventory STKVols > /dev/null 2> /dev/null
