# crontabs/silo_balance.fermi
13 3 * * * enstore /home/enstore/enstore/sbin/ecron silo-balance >/dev/null 2>/dev/null
