# crontabs/silo_balance.fermi
13 3 * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron silo-balance >/dev/null 2>/dev/null
