# crontabs/silo_balance.cdfen
13 3,11,19 * * * enstore /home/enstore/enstore/sbin/ecron silo-balance >/dev/null 2>/dev/null
