# crontabs/summary_plots
#The crontabs with two pound signs at the beginning of the comment are
# not known to be SDE ready.

# summary/overall tape burn rates
0 * * * * enstore source /usr/local/etc/setups.sh; setup enstore; $ENSTORE_DIR/sbin/ecron -p summary-burn-rate  $ENSTORE_DIR/src/summary_plots.py -t > /dev/null 2> /dev/null
