# crontabs/tab_flipping.d0en
00 16 * * * enstore /home/enstore/enstore/sbin/ecron -p write_tabs_audit /home/enstore/enstore/bin/write_tabs_audit > /dev/null 2> /dev/null
10 08 * * * enstore /home/enstore/enstore/sbin/ecron -p tab_flipping_watch /home/enstore/enstore/src/tab_flipping_nanny.py /diska/tape-inventory/TAB_FLIPPING_WATCH > /dev/null 2> /dev/null
20 08 * * * enstore /home/enstore/enstore/sbin/ecron -p tab_flipping_watch-samlto2 /home/enstore/enstore/src/tab_flipping_nanny.py -l samlto,samlto2 -o /diska/tape-inventory/TAB_FLIPPING_WATCH-samlto2 > /dev/null 2> /dev/null
