# crontabs/tab_flipping
0 * * * * enstore /home/enstore/enstore/sbin/ecron -p closing_tab_flipping_job /home/enstore/enstore/src/operation.py auto_close_all > /dev/null 2> /dev/null
