# sbin/adiclogs.cron
10,25,40,55 * * * * $ENSTORE_SBIN/ecron -p adiclogs adiclogs > /dev/null 2> /dev/null
18 * * * * $ENSTORE_SBIN/ecron -p adicmirror mirror $ENSTORE_DIR/sbin/adicmirror > /dev/null 2> /dev/null
