# sbin/backup.cron
10 * * * * $ENSTORE_SBIN/ecron -p backup ecmd backup > /dev/null 2> /dev/null
