# sbin/inqPlotUpdate.cron
0 1 * * * $ENSTORE_SBIN/ecron -p inqPlotUpdate -setup gnuplot -d /usr/hppc_home/berman/cron ecmd ic --plot --mc=grau.media_changer > /dev/null 2> /dev/null
30 1 * * * $ENSTORE_SBIN/ecron -p inqPlotUpdate -setup gnuplot -d /usr/hppc_home/berman/cron ecmd ic --plot > /dev/null 2> /dev/null
