# sbin/make_cron_plot_page.cron
0 3 * * * $ENSTORE_SBIN/ecron -p make_cron_plots make_cron_plots  > /dev/null 2> /dev/null
