# sbin/rdist-log.cron
45 * * * * $ENSTORE_SBIN/ecron -p rdist-log rdist-log  > /dev/null 2> /dev/null
